# watcher/mem-qdrant-watcher.service
[Unit]
Description=TrueRecall v1 - Real-Time Memory Watcher
After=network.target

[Service]
Type=simple
User=<USER>
WorkingDirectory=<INSTALL_PATH>/true-recall-v1/watcher
Environment="QDRANT_URL=http://<QDRANT_IP>:6333"
Environment="QDRANT_COLLECTION=memories_tr"
Environment="OLLAMA_URL=http://<OLLAMA_IP>:11434"
Environment="EMBEDDING_MODEL=snowflake-arctic-embed2"
Environment="USER_ID=<USER_ID>"
ExecStart=/usr/bin/python3 <INSTALL_PATH>/true-recall-v1/watcher/realtime_qdrant_watcher.py --daemon
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
